# inbox/cron.d/inbox-cleanup
*/3 * * * * service /cleanup.sh /maildir >/dev/null 2>&1
